# rootfsimg/inittab
::sysinit:/bin/busybox --install -s
/dev/console::sysinit:-/bin/bash /etc/init.d/rcS
/dev/console::sysinit:-/bin/bash
/dev/console::respawn:-/bin/bash