# rootfsimg/etc/inittab
::sysinit:/bin/busybox --install -s
/dev/console::sysinit:-/bin/sh /etc/init.d/rcS
/dev/console::respawn:-/bin/sh /etc/init.d/rcS